# nextcloud/nextcloud-cron.service
[Unit]
Description=vdsm hourly logrotate
BindsTo=nginx.service
[Service]
Type=oneshot
ExecStart=/usr/bin/php -f /usr/share/webapps/nextcloud/cron.php
User=http
